# service_files/cassie.conf
# Cassie daemon upstart config file

description "Cassie Bot"

start on runlevel [2345]
stop on runlevel [!2345]

exec /opt/cassie-bot/cassie.py -L INFO -c /opt/cassie-bot/cassie.conf -f
